# awesome/.config/autostart/clipit.desktop
[Desktop Entry]
Name=ClipIt
Exec=clipit
Terminal=false
Type=Application
