# autostart/.config/autostart/awdctl.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=awdctl
Exec=awdctl -d
Terminal=false
StartupNotify=false
